# Software/systemd/blackbox-web.service
[Unit]
Description=Holiday Blackbox Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Holiday-blackbox/Software
ExecStart=/usr/bin/python3 -m blackbox.web.app
Restart=on-failure

[Install]
WantedBy=multi-user.target

